# config/supervisor/conf.d/websockets.conf
[program:laravel-websockets]
process_name=%(program_name)s
directory=/var/www
command=php artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=www-app
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/websockets.log
stderr_logfile=/var/www/storage/logs/websockets-errors.log
